# services/schizophasia/deploy/registry/bootstrap/autokill.cron
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin:/usr/local/pgsql/bin/

*/1 * * * * /usr/local/bin/autokill.sh
